# frr-supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:zebra]
command=/usr/lib/frr/zebra
[program:bgpd]
command=/usr/lib/frr/bgpd
